# configs/expire.crontab
* * * * * /bin/bash -c 'curl --insecure -I https://10.5.161.169/expire >> /var/log/ssw/expire.log'
